# systemd/dev-mars-@{res}.path
[Unit]
Description=MARS auto-activation of primary device /dev/mars/@{res}
Documentation=https://github.com/schoebel/mars/docu/mars-manual.pdf
Requires=mars.service
After=mars.service

[Path]
PathExists=/dev/mars/@{res}
Unit=mnt-test-@escvar{res}.mount

[Install]
WantedBy=mars.service
