# systemd/mars-@{res}-trigger.path
[Unit]
Description=MARS auto-activation of primary device /dev/mars/@{res}
Documentation=https://github.com/schoebel/mars/docu/mars-manual.pdf
Requires=mars.service
After=mars.service

[Path]
PathChanged=/mars/resource-@{res}/systemd-want
PathModified=/mars/resource-@{res}/systemd-want
Unit=mars-trigger.service

[Install]
WantedBy=mars.service
