# systemd/mars-trigger.path
[Unit]
Description=MARS systemd trigger from /mars/userspace/systemd-trigger
Documentation=https://github.com/schoebel/mars/docu/mars-manual.pdf
Requires=mars.service
After=mars.service

[Path]
#PathExists=/mars/userspace/systemd-trigger
PathChanged=/mars/userspace/systemd-trigger
PathModified=/mars/userspace/systemd-trigger
Unit=mars-trigger.service

[Install]
WantedBy=mars.service
