# systemd/mars.path
[Unit]
Description=MARS block layer HA (activation by appearance of /mars/uuid)
Documentation=https://github.com/schoebel/mars/docu/mars-manual.pdf

[Path]
PathExists=/mars/uuid
Unit=mars.service
