# systemd-testing/mars-emergency.path
[Unit]
Description=MARS trigger from /mars/userspace/emergency-trigger
Documentation=https://github.com/schoebel/mars/docu/mars-user-manual.pdf
Requires=mars.service
After=mars.service

[Path]
#PathExists=/mars/userspace/emergency-trigger
PathChanged=/mars/userspace/emergency-trigger
PathModified=/mars/userspace/emergency-trigger
Unit=mars-emergency.service

[Install]
WantedBy=mars.service
