# systemd-testing/mars-emergency.service
[Unit]
Description=MARS emergency trigger service from /mars/userspace/emergency-trigger
Documentation=https://github.com/schoebel/mars/docu/mars-user-manual.pdf

[Service]
Type=oneshot
RemainAfterExit=no
# Important for avoiding failures: do not send any signals
KillMode=none
ExecStart=/usr/bin/marsadm cron

[Install]
WantedBy=mars.service
